# build/etc/service/nginx/run
#!/usr/bin/env bash

nginx -c /etc/nginx/current.conf